# metalcloud/netboot/metalagent.service
[Unit]
Description=Run agent

[Service]
ExecStart=/sbin/metalagent

[Install]
WantedBy=multi-user.target